# scripts/service/linux/0ctx-daemon.service
[Unit]
Description=0ctx Connector Runtime — cloud bridge + managed daemon lifecycle
Documentation=https://github.com/0ctx-com/0ctx
After=network.target

[Service]
Type=simple
ExecStart=%NODE_PATH% %CLI_ENTRY% connector run --quiet
Restart=on-failure
RestartSec=5s
# Escalating restart delays
StartLimitBurst=5
StartLimitIntervalSec=60

# Run as the installing user (see install instructions for user-level service)
# Environment
Environment=NODE_ENV=production

# Logging goes to journald automatically (use: journalctl --user -u 0ctx-daemon)

[Install]
WantedBy=default.target
